# roles/mgrote_youtubedl/templates/youtubedl.timer.j2
{{ file_header | default () }}
[Unit]
Description=Timer for youtubedl.

[Timer]
OnCalendar={{ ytdl_timer }}
RandomizedDelaySec=30 min

[Install]
WantedBy=timers.target multi-user.target
