# roles/mgrote_zfs_sanoid/templates/syncoid.timer.j2
{{ file_header | default () }}
[Unit]
Description=Timer for syncoid.

[Timer]
OnCalendar={{ sanoid_syncoid_timer }}

[Install]
WantedBy=timers.target multi-user.target zfs.target
