# roles/_mgrote.musterrolle_systemd/templates/rclone.timer.j2
{{ file_header | default () }}
[Unit]
Description=Timer for rclone syncs.

[Timer]
OnCalendar={{ item.timer }}
RandomizedDelaySec=30 min

[Install]
WantedBy=timers.target multi-user.target
